# server/publish-server.service
[Unit]
Description=gszsyy GitHub Pages project publish backend
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/data_1/wangwei/duomianma-sdk-weekly-pages
EnvironmentFile=/home/data_1/wangwei/duomianma-sdk-weekly-pages/server/.env
ExecStart=/usr/bin/env node /home/data_1/wangwei/duomianma-sdk-weekly-pages/server/publish-server.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
